# supervisord.conf
[supervisord]
nodaemon=true

[program:youtrack]
command=/bin/bash -c "/usr/local/youtrack/youtrack.sh start"